# .claude/hooks/com.knowledgeos.daily-autopilot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.knowledgeos.daily-autopilot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/jwpark95/Documents/Knowledge-OS/.claude/hooks/daily-autopilot.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/jwpark95/Documents/Knowledge-OS/.claude/logs/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jwpark95/Documents/Knowledge-OS/.claude/logs/launchd-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/jwpark95</string>
    </dict>
</dict>
</plist>
